# examples/fr_covidata_engine_pky_development.cron
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

# Use appointments to create fake lab results every minute
# Load fake lab results into the PKY survey and serial projects every two minutes
* * * * *   root /usr/bin/docker run --rm -it --env-file /ctsit/fr_covidata_engine/pky_dev.env fr_covidata_engine Rscript load_fake_data_into_upload.R
* * * * *   root /usr/bin/docker run --rm -it --env-file /ctsit/fr_covidata_engine/pky_dev.env fr_covidata_engine Rscript load_fake_data_from_serial_into_upload.R
*/2 * * * * root /usr/bin/docker run --rm -it --env-file /ctsit/fr_covidata_engine/pky_dev.env fr_covidata_engine Rscript load_results_into_pky_projects.R
